# deploy/cbr-chatbot.service
[Unit]
Description=CBR Customer Support Chatbot (FastAPI + Uvicorn)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/cbr-chatbot
Environment="CBR_DB_PATH=/opt/cbr-chatbot/data/cases.db"
EnvironmentFile=-/opt/cbr-chatbot/.env
ExecStart=/opt/cbr-chatbot/.venv/bin/uvicorn api.main:app --host 0.0.0.0 --port 8000 --workers 2
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
